# backup-configs.timer
[Unit]
Description=Timer for backupConfigs

[Timer]
OnCalendar=monthly
Persistent=yes
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target